# hello_web.service
[Unit]
Description=runs a hello world webapp
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/user/.volta/bin/node /var/www/src/index.js
User=username
Group=user
Restart=always
RestartSec=10
TimeoutStopSec=90
SyslogIdentifier=hello_app

[Install]
WantedBy=multi-user.target
